# build/super/template.conf
[supervisord]
nodaemon=true

[program:PROGRAM_NAME]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/bin/php /var/www/artisan queue:work QUEUE_CONNECTION --queue=QUEUES --timeout=TIMEOUT --max-time=MAX_TIME --once -v
autostart=true
autorestart=true
user=root
numprocs=PROCESS_NUM
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/super/LOG_FILE.log
